# gunicorn.conf
description "Gunicorn application server handling myproject"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid rwholey 
setgid www-data
chdir /home/rwholey/bots/jeopardy_bot/trabek_bot

exec venv/bin/gunicorn --workers 3 --bind unix:/home/rwholey/bots/jeopardy_bot/trabek_bot/trabek_bot.sock trabek_bot.wsgi:application
